# deploy/systemd/meatbot.service
[Unit]
Description=Meat Plant Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=meatbot
Group=meatbot
WorkingDirectory=/opt/BOT_Miasokombinat
EnvironmentFile=/opt/BOT_Miasokombinat/.env
ExecStart=/opt/BOT_Miasokombinat/venv/bin/python main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
