# deploy/systemd/meatbot-backup.service
[Unit]
Description=Backup SQLite database for Meat Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=meatbot
Group=meatbot
WorkingDirectory=/opt/BOT_Miasokombinat
EnvironmentFile=/opt/BOT_Miasokombinat/.env
ExecStart=/opt/BOT_Miasokombinat/venv/bin/python scripts/backup_db.py
